# prometheus-quit.service.in
[Unit]
Description=Terminate Prometheus Boot Screen
After=rc-local.service systemd-user-sessions.service
Before=getty@tty1.service

[Service]
ExecStart=-@CMAKE_INSTALL_PREFIX@/sbin/prometheus quit 
Type=oneshot
TimeoutSec=20
