# systemd/prometheus-poweroff.service.in
[Unit]
Description=Show Prometheus Power Off Screen
After=getty@tty1.service 
Before=poweroff.service
DefaultDependencies=no
ConditionKernelCommandLine=!prometheus.enable=0

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/sbin/prometheusd --stage shutdown
Type=forking
